# synop.cron
PATH=/usr/bin:/bin:/opt/eccodes/bin
*/10 * * * * cd /usr/src/app && python3 /usr/local/bin/flask load_observations > /proc/1/fd/1 2>/proc/1/fd/2
